# data/services/systemd/swayosd-libinput-backend.service.in
[Unit]
Description=SwayOSD LibInput backend for listening to certain keys like CapsLock, ScrollLock, VolumeUp, etc...
Documentation=https://github.com/ErikReider/SwayOSD
PartOf=graphical.target
After=graphical.target

[Service]
Type=dbus
BusName=org.erikreider.swayosd
ExecStart=@bindir@/swayosd-libinput-backend
Restart=on-failure

[Install]
WantedBy=graphical.target
